# supervisor/single_consumer.conf
[program:rabbitmq_scaling0]
command=php -dmemory_limit=512M /app/bin/console rabbitmq:consumer -m 1000 scaling0 --env=dev
process_name=rabbitmq_scaling0_%(process_num)02d
numprocs=1
directory=/tmp
autostart=true
autorestart=true
startsecs=3
startretries=10
user=www-data
redirect_stderr=false
stdout_logfile=/app/var/log/supervisor.rabbitmq.scaling0.out.log
stdout_capture_maxbytes=1MB
stderr_logfile=/app/var/log/supervisor.rabbitmq.scaling0.error.log
stderr_capture_maxbytes=1MB

[program:kafka_scaling0]
command=php -dmemory_limit=512M /app/bin/console kafka:consumer:run scaling0 --env=dev
process_name=kafka_scaling0_%(process_num)02d
numprocs=1
directory=/tmp
autostart=true
autorestart=true
startsecs=3
startretries=10
user=www-data
redirect_stderr=false
stdout_logfile=/app/var/log/supervisor.kafka.scaling0.out.log
stdout_capture_maxbytes=1MB
stderr_logfile=/app/var/log/supervisor.kafka.scaling0.error.log
stderr_capture_maxbytes=1MB
